# debian/wlanpi-webui.socket
# wlanpi-webui.socket

[Unit]
Description=wlanpi-webui gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
User=wlanpi

[Install]
WantedBy=sockets.target
